# external/hls/2019.08-bcolbf/etc/profile.d/dependencies-setup.csh
test $?GCC_ROOT != 0 || source /cvmfs/cms.cern.ch/slc7_amd64_gcc820/external/gcc/8.2.0-bcolbf/etc/profile.d/init.csh
test $?GMAKE_ROOT != 0 || source /cvmfs/cms.cern.ch/slc7_amd64_gcc820/external/gmake/4.3/etc/profile.d/init.csh
